# configs/skiff/standard/root_overlay/etc/systemd/system/systemd-modules-load.service.d/after-init.conf
[Unit]
After=usr-lib-modules.mount
Requires=usr-lib-modules.mount
