# configs/skiff/standard/root_overlay/etc/systemd/system/systemd-modules-load.service.d/after-init.conf
[Unit]
After=skiff-init.service
Requires=skiff-init.service
