# configs/skiff/core/root_overlay/etc/systemd/system/skiff-core.service
[Unit]
Description=SkiffOS core setup
Requires=docker.service
After=docker.service
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/opt/skiff/scripts/skiff-core.sh
Restart=on-failure
RestartSec=10


