# configs/pine64/phone/root_overlay/etc/udev/rules.d/60-flashlight.rules
SUBSYSTEM=="leds", KERNEL=="white:flash", ACTION=="add", RUN+="/bin/chmod -R a+w /sys%p"