# configs/browser/v86/root_overlay/etc/systemd/system/systemd-journal-flush.service.d/01-browser-v86-journal-flush.conf
[Unit]
# Delay flushing till system is done booting.
Before=

[Service]
TimeoutSec=0
