# docker/docker.service
[Unit]
Description=Placeholder for docker service.

[Service]
Type=oneshot
ExecStart=/usr/bin/env
RemainAfterExit=true
