# configs/util/adbd/root_overlay/usr/lib/systemd/system/adbd.service
[Unit]
Description=ADB Daemon

[Service]
Restart=always
ExecStart=/opt/skiff/scripts/start-adbd.sh

[Install]
WantedBy=multi-user.target