# configs/skiff/standard/root_overlay/etc/systemd/system/docker.service.d/after-init.conf
[Unit]
After=skiff-init.service
Requires=skiff-init.service
After=mnt-persist.mount
Requires=mnt-persist.mount