# configs/skiff/standard/root_overlay/etc/systemd/system/docker.service.d/ramdisk.conf
[Service]
Environment=DOCKER_RAMDISK=true
