# configs/skiff/standard/root_overlay/usr/lib/systemd/system/persist-resize.service
[Unit]
Description=SkiffOS resize persist partition
DefaultDependencies=no
Before=sysinit.target

[Service]
Type=oneshot
ExecStart=/opt/skiff/scripts/resize-persist.sh
RemainAfterExit=true
