# configs/artik/710/root_overlay/etc/udev/rules.d/10-local.rules
#ACTION=="add", KERNEL=="hci0", RUN+="/etc/bluetooth/hciconf.sh"
