# configs/artik/710/root_overlay/etc/udev/rules.d/99-adb-restart.rules
ACTION=="change", SUBSYSTEM=="udc", RUN+="/bin/sh -c '/bin/echo c0040000.dwc2otg > /sys/kernel/config/usb_gadget/g1/UDC'"
