# configs/skiff/standard/root_overlay/etc/systemd/system/skiff-init.service
[Unit]
Description=SkiffOS early-init

[Service]
Type=oneshot
ExecStart=/opt/skiff/scripts/mount-all.sh
