# configs/skiff/standard/root_overlay/etc/systemd/system/swapon.service
[Unit]
Description=Turn on swap

[Service]
Type=oneshot
ExecStart=/sbin/swapon LABEL=swap

[Install]
WantedBy=swap.target
