# configs/apps/balena/root_overlay/usr/lib/systemd/system/balena-engine.service.d/ramdisk.conf
[Service]
Environment=DOCKER_RAMDISK=true
Environment=BALENA_ENGINE_RAMDISK=true
