# configs/virt/qemu/root_overlay/etc/systemd/system/sshd.service.d/01-virt-qemu-startup-timeout.conf
[Service]
# Give longer for ssh-keygen to work.
TimeoutStartSec=600
